# rpi5 services/ngrok-arduino.service
[Unit]
Description=ngrok tunnel for Arduino UI
After=network-online.target arduino-web.service
Wants=network-online.target arduino-web.service

[Service]
Environment=NGROK_AUTHTOKEN=<your_token>
ExecStart=/usr/local/bin/ngrok http 5000 --log stdout
WorkingDirectory=/home/fcp1
Restart=always
User=fcp1
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
